# root/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash

umask ${UMASK_SET:-022}

exec /usr/sbin/php-fpm7 -F
